# supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log;
pidfile=/var/run/supervisord.pid;
loglevel=debug

[program:kms]
command=/bin/bash /kms.sh
redirect_stderr=true

[program:openvidu-server]
command=/bin/bash -c "java -jar -Dserver.ssl.key-store=/opt/openvidu/keystore.jks -Dserver.ssl.key-store-password=$KEYSTORE_PW -Dserver.ssl.key-alias=$KEYSTORE_ALIAS openvidu-server.jar"
redirect_stderr=true
